# recipes-bsp/advantech-initscripts/advantech-initscripts/bootcount.service
[Unit]
Description=Bootcount service
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/lib/bootcount.sh

[Install]
WantedBy=custom.target
